# root/etc/services.d/php-fpm/run
#!/bin/bash
exec /usr/bin/php-fpm -y /config/nginx/php-fpm.conf -F
